# linux-monitor.service
[Unit]
Description=Linux System Monitor with MQTT and Home Assistant Discovery (Python)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/linux_mqtt_ha
ExecStart=/usr/bin/python3 /opt/linux_mqtt_ha/mqtt_linux_monitoring.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables (optional - can be set in the script instead)
# Environment=MQTT_BROKER=localhost
# Environment=MQTT_PORT=1883

[Install]
WantedBy=multi-user.target
